# recipes-bsp/tools/resize-helper/resize-helper.sh.in
#! /bin/sh
### BEGIN INIT INFO
# Provides: Resize root filesystem to fit available disk space
# Description: Start resize helper, then remove own init from runlevels,
# as it is required only for first boot. run settle to provide partition
# table details.
### END INIT INFO

DESC="resize helper"

case $1 in
	start)
		echo "Starting $DESC"
		@sbindir@/resize-helper
	;;
  *)
		echo "Usage: @sysconfdir@/init.d/resize-helper.sh {start}" >&2
		exit 1
	;;
esac

exit 0

# vim:noet
