# recipes-graphics/wayland/weston-init/73-pulseaudio-hdmi.rules
SUBSYSTEMS=="drm", ACTION=="change", RUN+="/usr/bin/pulseaudio_hdmi_switch.sh"
